# Code red alert task/[root]/usr/lib/systemd/user/update-notifier-crash.path
[Unit]
Description=Path trigger for Apport crash notifications
PartOf=graphical-session.target
ConditionPathExists=!/var/lib/apport/autoreport

[Path]
PathChanged=/var/crash/
