# Code red alert task/[root]/etc/systemd/system/cron.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/sbin/cron -f -L 0
